# backend/sigma_rules/persistence/t1547.001_run_key_modified.yml
title: Run / RunOnce registry persistence
id: 0a4d1f23-c001-4ec0-a5b2-4e1f5b7fb105
status: experimental
description: |
  Detects writes to autorun registry keys. While legitimate
  installers also write here, the path is a top-three persistence
  mechanism on Windows; combined with M11.b anomaly detection
  (parent-process baselining), this rule fires on attacker-driven
  writes specifically.
references:
  - https://attack.mitre.org/techniques/T1547/001/
author: edr-curated
date: 2026-05-10
tags:
  - attack.persistence
  - attack.t1547.001
logsource:
  category: registry_set
  product: windows
detection:
  selection:
    registry.path|contains:
      - '\Software\Microsoft\Windows\CurrentVersion\Run\'
      - '\Software\Microsoft\Windows\CurrentVersion\RunOnce\'
      - '\Software\Wow6432Node\Microsoft\Windows\CurrentVersion\Run\'
  condition: selection
fields:
  - registry.path
  - registry.value
  - process.executable
falsepositives:
  - Software installers.
  - Operator-driven autorun configuration.
level: medium
